# warema_bridge/rootfs/etc/s6-overlay/s6-rc.d/warema-bridge/run
#!/usr/bin/execlineb -P
cd /srv
exec node bridge.js
